# service/megapi-ap.service
[Unit]
Description=As an AP for connection

[Service]
Type=idle
ExecStart=/usr/bin/create_ap wlan0 eth0 megapi

[Install]
WantedBy=multi-user.target
